# etc/service/dovecot/run
#!/bin/bash -e

exec /usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf 2>&1
