# etc/service/sshd/run
#!/bin/bash -e

mkdir -p /var/run/sshd

exec /usr/sbin/sshd -D -e 2>&1
