# etc/service/postgrey/run
#!/bin/bash -e

# Logging goes to syslog (mail).

exec /usr/sbin/postgrey --inet=10023 2>&1
