# etc/cron.daily/clamav-freshclam
#!/bin/bash -e

/usr/bin/freshclam --quiet -l /var/log/clamav/freshclam.log
