# etc/service/clamav-freshclam/run
#!/bin/bash -e

exec chpst -u clamav:clamav /usr/bin/freshclam 2>&1
